# device0/cfg/etc/systemd/system/podman-auto-update.timer.d/override.conf
[Timer]
OnCalendar=
RandomizedDelaySec=0
OnBootSec=15
OnUnitActiveSec=15
